# com.feishu.claude.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.feishu.claude.bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/FULL/PATH/TO/python3</string>
        <string>/FULL/PATH/TO/bot.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/FULL/PATH/TO/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/FULL/PATH/TO/bot.log</string>
    <key>WorkingDirectory</key>
    <string>/FULL/PATH/TO/</string>
</dict>
</plist>
